# scanew-web.service
[Unit]
Description=Scannew Web Application
After=network.target

[Service]
# IMPORTANT: Replace with the username you are running the app as on your server
User=outjet
# IMPORTANT: Replace with the group for that user
Group=outjet

# IMPORTANT: Replace with the absolute path to the project directory on your server
WorkingDirectory=/home/outjet/scanew

# IMPORTANT: Replace with the absolute path to the gunicorn executable in your virtual environment
ExecStart=/home/outjet/scanew/myenv/bin/gunicorn --workers 4 --bind 0.0.0.0:8000 "src.web_app:create_app()"

# Environment variables for Flask and extensions
# IMPORTANT: Set a strong, unique SECRET_KEY
Environment="SECRET_KEY=your_very_secret_key_here"
# IMPORTANT: Set your Google OAuth credentials
Environment="GOOGLE_OAUTH_CLIENT_ID=your_google_oauth_client_id"
Environment="GOOGLE_OAUTH_CLIENT_SECRET=your_google_oauth_client_secret"
# IMPORTANT: Set your OpenAI API key
Environment="OPENAI_API_KEY=your_openai_api_key"
# IMPORTANT: Set your Redis URL if it's not the default
Environment="REDIS_URL=redis://localhost:6379"

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
